# etc/cron.d/proximus
# reload noauth squid config

*/2 * * * * root test -x /usr/share/proximus/proximus_noauth.php && /usr/share/proximus/proximus_noauth.php
