# etc/cron.d/proximus
# reload noauth squid config

*/2 * * * * root test -x /usr/share/proximus/proximus_noauth_reload.sh && /usr/share/proximus/proximus_noauth_reload.sh
