# proximus-redirector/debian/postinst
#! /bin/sh
# postinst script for proximus
#
# see: dh_installdeb(1)

set -e

echo Reloading Squid to use new Proximus package
invoke-rc.d squid reload

